# dummy-api.service
[Unit]
Description=Dummy API built in Kotlin using Spring Boot

[Service]
ExecStart=/var/dummy-api/dummy-api.jar
Restart=on-failure
User=root
RestartSec=1

[Install]
WantedBy=multi-user.target